# supervisor/logstash-indexer.conf
[program:logstash-indexer]
directory=/
command=/opt/logstash/bin/logstash agent -f /etc/logstash/logstash-indexer.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/logstash-indexer.out.log
stderr_logfile=/var/log/supervisor/logstash-indexer.err.log
